# services/led-matrix-web.service
[Unit]
Description=LED Matrix Web Control Panel
Documentation=https://github.com/RynAgain/LED_MATRIX-Project
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/LED_MATRIX-Project
ExecStart=/home/pi/LED_MATRIX-Project/venv/bin/python3 -m src.web.app
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-matrix-web

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/LED_MATRIX-Project

[Install]
WantedBy=multi-user.target
