# services/led-matrix-updater.service
[Unit]
Description=LED Matrix Auto-Updater Service
Documentation=https://github.com/RynAgain/LED_MATRIX-Project
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ryn
Group=ryn
WorkingDirectory=/home/ryn/LED_MATRIX-Project
ExecStart=/bin/bash /home/ryn/LED_MATRIX-Project/scripts/update.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-matrix-updater

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ryn/LED_MATRIX-Project
Environment=HOME=/home/ryn

[Install]
WantedBy=multi-user.target
