# services/led-matrix.service
[Unit]
Description=LED Matrix Display Service
Documentation=https://github.com/RynAgain/LED_MATRIX-Project
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Must run as root for GPIO/hardware access (rpi-rgb-led-matrix requirement)
User=root
Group=root
WorkingDirectory=/home/pi/LED_MATRIX-Project
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/LED_MATRIX-Project/venv/bin/python3 /home/pi/LED_MATRIX-Project/src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=led-matrix

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/LED_MATRIX-Project

# Hardening (ProtectHome removed - project lives under /home and needs write access)
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
